# prototype/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
user=app

# MediaMTX media server — RTMP/RTSP/HLS/WebRTC ingest + playback + recording.
[program:mediamtx]
command=/usr/local/bin/mediamtx /mediamtx.yml
autorestart=true
startsecs=3
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Go JSON API — proxies MediaMTX over localhost (same container network namespace).
# Backend defaults already point at localhost:{8888,8889,9997,9996}; no env needed.
[program:backend]
command=/usr/local/bin/mediamtx-console
autorestart=true
startsecs=3
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
